# supervisord/mbus_poll.conf
[program:mbus_poll]
command=python2.7 /home/pi/bin/mbus_poll.py
directory=/home/pi/
user=pi
autostart=true
autorestart=false
